# haproxy/supervisord-incron.conf
[program:incron]
command=/usr/sbin/incrond
numprocs=1
autostart=true
autorestart=false
startsecs=0
